# ubuntu/profile.d/x-show-minit-banner.sh
if [ -f /etc/banner.minit.txt ]; then
  cat /etc/banner.minit.txt
fi